# rootfs/etc/s6-overlay/s6-rc.d/nfs-serve/run
#!/usr/bin/env sh

exec rpc.mountd -F -N 2 -N 3
