# rootfs/etc/s6-overlay/s6-rc.d/nfs-configure/up
#!/usr/bin/env sh

echo "user_allow_other" >> /etc/fuse.conf

exportfs -r

rpc.nfsd -N 3
